# confs/supervisord.conf
[supervisord]
nodaemon=true

[program:fastapi]
command=python3 app.py

[program:celery]
command=celery -A tasks worker --loglevel=info